# src/bios-agent-cm.service.in
[Unit]
Description=BIOS computation services (min, max, averages)
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/bios/bios-db-rw
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/bios-agent-cm
Restart=allways

[Install]
WantedBy=bios.target
