# src/fty-metric-compute.service.in
[Unit]
Description=42ity computation services (min, max, averages)
Requires=malamute.service
After=malamute.service
PartOf=fty.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/bios/etc/default/bios
EnvironmentFile=-/usr/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios
EnvironmentFile=-/etc/default/bios__%n.conf
EnvironmentFile=-/etc/default/bios-db-rw
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-metric-compute
Restart=always

[Install]
WantedBy=bios.target
