# src/fty-metric-compute.service.in
[Unit]
Description=42ity computation services (min, max, averages)
Requires=malamute.service
After=malamute.service
PartOf=bios.target

[Service]
Type=simple
User=bios
EnvironmentFile=-/usr/share/fty/etc/default/fty
EnvironmentFile=-/usr/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/fty
EnvironmentFile=-/etc/default/fty__%n.conf
EnvironmentFile=-/etc/default/fty-db-rw
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/fty-metric-compute
Restart=always

[Install]
WantedBy=bios.target
